# meta-iris/recipes-core/iris-utils/files/fstrim.cron
#!/bin/sh
LOG=/var/log/fstrim.log
echo "Performing fstrim: $(date -R)" >> $LOG
time /sbin/fstrim -v /data >> $LOG 2<&1
